# DEBIAN/postinst
#!/bin/sh

cd /Library/LaunchDaemons
launchctl load com.apple.sharingd.plist